# rslsync-backup.timer
[Unit]
Description=Perform backup of selected btsync directories

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
